# srv.service
[Unit]
Description=Bookmark Manager
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/bookmark-manager
ExecStart=/home/exedev/bookmark-manager/bookmarks -listen :8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
